# etc/systemd/system/vision.service
[Unit]
After=multi-user.target.service

[Service]
ExecStart=/usr/local/bin/vision-script

[Install]
WantedBy=default.target
